# deploy/supervisor/wedding-website.conf
[program:wedding-website]
command=[path to gunicorn start script]
user=mariavelarde
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=[path to log file]
